# host/extras/system/usb-das-keepalive.target
[Unit]
Description=Enable keepalive timers for discovered USB-DAS members
