# systemd/proxnix-gc.timer
[Unit]
Description=Periodically garbage-collect stale proxnix stage directories

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
