# host/runtime/systemd/proxnix-gc.timer
[Unit]
Description=Periodically garbage-collect stale proxnix host state

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
