# host/runtime/systemd/proxnix-reconcile.service
[Unit]
Description=Reconcile proxnix managed NixOS LXCs
After=network-online.target pve-cluster.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/proxnix-reconcile
