# host/runtime/systemd/proxnix-gc.service
[Unit]
Description=Garbage-collect stale proxnix host state
Documentation=https://github.com/your-org/proxnix

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/proxnix-host gc
