# host/runtime/systemd/proxnix-reconcile@.service
[Unit]
Description=Reconcile proxnix managed NixOS LXC %i
After=network-online.target pve-cluster.service
Wants=network-online.target

[Service]
Type=oneshot
# LXC pre-start triggers this asynchronously before the guest can accept
# pct exec. Give the start path time to reach the running state.
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/sbin/proxnix-reconcile --vmid %i
