# host/runtime/systemd/proxnix-cache-reconcile.timer
[Unit]
Description=Retry pending proxnix shared-cache uploads

[Timer]
OnBootSec=10min
OnUnitInactiveSec=15min
RandomizedDelaySec=2min
Persistent=true

[Install]
WantedBy=timers.target
