# host/runtime/systemd/proxnix-gc.timer
[Unit]
Description=Periodically garbage-collect stale proxnix host state

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
