# host/runtime/systemd/proxnix-reconcile@.service
[Unit]
Description=Reconcile proxnix managed NixOS LXC %i
After=network-online.target pve-cluster.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/proxnix-host reconcile --vmid %i
