# host/runtime/systemd/proxnix-reconcile.service
[Unit]
Description=Run automatic proxnix host event reconciliation
After=network-online.target pve-cluster.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/proxnix-host reconcile --auto-tag
