# system/cluster-folder-git-reconcile.service
[Unit]
Description=cluster-folder-git periodic reconcile
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cluster-folder-git reconcile-all
